# driver-kafka-on-es/deploy/templates/range-server.service
[Unit]
Description=Range Server
After=network.target

[Service]
ExecStart={{ rs_dir }}/target/release/range-server --config ./rs-config.yaml --log ./rs-log-config.yaml
WorkingDirectory={{ rs_dir }}
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
