# driver-kafka-on-es/deploy/templates/range-server.service
[Unit]
Description=Range Server
After=network.target

[Service]
ExecStart=/usr/local/bin/range-server start --config ./rs-config.yaml --log ./rs-log-config.yaml
Environment='ES_ADDR=0.0.0.0:10911'
Environment='ES_ADVERTISE_ADDR={{ private_ip }}:10911'
Environment='ES_PD={{ pd_0_ip }}:12378'
WorkingDirectory={{ rs_dir }}
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
